# petdoor.service
[Unit]
Description=Pet Door Controller
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/petdoor/raspberry_pi
ExecStart=/usr/bin/python3 /home/pi/petdoor/raspberry_pi/main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
